# tool/systemd/coin_trade.service
[Unit]
Description=COIN TRADE service.

[Service]
ExecStart=/bin/bash /usr/bin/coin_trade.sh

[Install]
WantedBy=multi-user.target
